# docker/crontab
0 */12 * * * php /app/index.php >> /var/log/script.log
